# scripts/painkiller.hillshade.hammock.service
[Unit]
Description=painkiller/hillshade

[Service]
ExecStart=/usr/bin/python /home/arctair/apps/painkiller/hillshade/main.py 0
Type=simple
WorkingDirectory=/home/arctair/apps/painkiller/hillshade

[Install]
WantedBy=default.target
